# systemd/nikko-music-sync.service
[Unit]
Description=NikkoMusicHub Dropbox Sync
After=network.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/srv/nikko-music/app
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/nikko-music/venv/bin/python -m app.services.sync_runner
