# app/systemd/nikko-music-watchdog.service
[Unit]
Description=NikkoMusicHub Watchdog Health Check
After=network.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/srv/nikko-music/app
EnvironmentFile=-/srv/nikko-music/data/nikko.env
Environment="PATH=/srv/nikko-music/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/srv/nikko-music/venv/bin/python -m app.services.health_check
StandardOutput=journal
StandardError=journal
