# app/systemd/nikko-music-player.service
[Unit]
Description=NikkoMusicHub mpv Player
After=network.target sound.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/srv/nikko-music/app
Environment="PATH=/srv/nikko-music/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
ExecCondition=/srv/nikko-music/venv/bin/python -m app.services.mpv_check
ExecStart=/usr/bin/mpv --no-video --shuffle --loop-playlist=inf --input-ipc-server=/tmp/nikko-mpv.sock --playlist=/srv/nikko-music/playlist.m3u --log-file=/srv/nikko-music/logs/player.log
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
