# cloud/systemd/nikko-music-hub-cloud.service
[Unit]
Description=NikkoMusicHub Cloud Management Platform
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/nikko-cloud/cloud
Environment=PYTHONUNBUFFERED=1
Environment=NIKKO_CLOUD_SECRET_KEY=%H-%M-%S-nikko-cloud
ExecStart=/opt/nikko-cloud/cloud/venv/bin/uvicorn cloud.app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
